# deploy/systemd/user/garoon-icloud-sync@.timer
[Unit]
Description=Run garoon-icloud-sync for profile %i every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true
Unit=garoon-icloud-sync@%i.service

[Install]
WantedBy=timers.target
